# files/gitlab-misc-backup.timer
[Unit]
Description=Backup gitlab misc

[Timer]
OnCalendar=*-*-* 2,14:00
Persistent=true

[Install]
WantedBy=timers.target
