# files/gitlab-remove-old-backup.service
[Unit]
Description=Remove old backup of gitlab

[Service]
Type=oneshot
ExecStart=/etc/gitlab/gitlab-remove-old-backup.sh
