# files/gitlab-remove-old-backup.timer
[Unit]
Description=Remove old backup of gitlab

[Timer]
OnCalendar=*-*-* 2,14:00

[Install]
WantedBy=timers.target
